# images/ochograph/resources/supervisor/ochograph.conf
[program:ochograph]
command=python /opt/ochograph/pod/pod.py
